# deploy/coded-api-backup.timer
[Unit]
Description=Daily Coded API backup

[Timer]
OnCalendar=*-*-* 03:17:00
Persistent=true
RandomizedDelaySec=10m

[Install]
WantedBy=timers.target
